# easyconfig/files/etc/init.d/easyconfig_modemdetect
#!/bin/sh /etc/rc.common

START=91
STOP=91

start() {
	[ -e /sys/kernel/debug/usb/devices ] || exit 0
	(/usr/share/easyconfig/modem/detect.sh >/dev/null) &
}
